# protected/data/crontab/dev.all4holidays
0 */1 * * * /usr/bin/indexer --config /var/www/dev.all4holidays/sphinx/dev/sphinx.conf --rotate a4h_deal > /dev/null 2>&1
*/5 * * * * /usr/bin/indexer --config /var/www/dev.all4holidays/sphinx/dev/sphinx.conf --rotate a4h_deal_delta > /dev/null 2>&1